# AppsDocker/supervisor/supervisord-ratchet.conf
[program:ratchet]
command=/start-ratchet.sh
numprocs=1
autostart=true
autorestart=true