# supervisor/supervisord-webwolf.conf
[program:webwolf]
command=/runfiles/start-webwolf.sh
numprocs=1
autostart=true
autorestart=true