# supervisor/supervisord-webgoat.conf
[program:webgoat]
command=/start-webgoat.sh
numprocs=1
autostart=true
autorestart=true