# supervisor/supervisord-djangonv.conf
[program:djangoNV]
command=/runfiles/start-djangonv.sh
numprocs=1
autostart=true
autorestart=true