# supervisor/supervisord-dvws-node.conf
[program:juiceshop]
command=/start-dvws-node.sh
numprocs=1
autostart=true
autorestart=true