# supervisor/supervisord-juiceshop.conf
[program:juiceshop]
command=/start-juiceshop.sh
numprocs=1
autostart=true
autorestart=true