# supervisor/supervisord-mysqld.conf
[program:mysqld]
command=/runfiles/start-mysqld.sh
numprocs=1
autostart=true
autorestart=true