# supervisor/supervisord-mailcatcher.conf
[program:mailcatcher]
command=/start-mailcatcher.sh
numprocs=1
autostart=true
autorestart=true