# supervisor/supervisord-wrongsecrets.conf
[program:wrongsecrets]
command=/runfiles/start-wrongsecrets.sh
numprocs=1
autostart=true
autorestart=true