# supervisor/supervisord-wrongsecrets.conf
[program:webgoat]
command=/start-wrongsecrets.sh
numprocs=1
autostart=true
autorestart=true