# supervisor/supervisord-railsgoat.conf
[program:railsgoat]
command=/start-railsgoat.sh
numprocs=1
autostart=true
autorestart=true