# supervisor/supervisord-dvgql.conf
[program:dvgql]
command=/runfiles/start-dvgql.sh
numprocs=1
autostart=true
autorestart=true