# supervisor/supervisor-wrongsecrets.conf
[program:wrongsecrets]
command=/start-wrongsecrets.sh
numprocs=1
autostart=true
autorestart=true