# supervisor/supervisord-dvws.conf
[program:juiceshop]
command=/start-dvws.sh
numprocs=1
autostart=true
autorestart=true