# supervisor/supervisord-juiceshop.conf
[program:juiceshop]
command=/runfiles/start-juiceshop.sh
numprocs=1
autostart=true
autorestart=true